# node-red/rootfs/etc/services.d/nodered/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Take down the S6 supervision tree when Node-RED fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Node-RED crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Node-RED stoped, restarting..."
