# node-red/rootfs/etc/services.d/nodered/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Starts Node-RED
# ==============================================================================
export NODE_PATH=/opt/node_modules:/data/node_modules
declare -a options

options+=(--settings "/etc/node-red/config.js")

cd /opt || bashio::exit.nok "Could not change directory to Node-RED"

bashio::log.info "Starting Node-RED..."

exec npm start -- "${options[@]}"
