# node-red/rootfs/etc/s6-overlay/s6-rc.d/nodered/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Starts Node-RED
# ==============================================================================
export NODE_PATH=/opt/node_modules:/data/node_modules
declare -a options
declare max_old_space_size

if bashio::config.exists 'max_old_space_size'; then
    max_old_space_size=$(bashio::config 'max_old_space_size')
    bashio::log.info "Setting max_old_space_size to ${max_old_space_size}"
    export NODE_OPTIONS="--max_old_space_size=${max_old_space_size}"
fi

if bashio::config.true 'safe_mode'; then
    bashio::log.info "Setting safe mode startup"
    options+=(--safe)
fi

options+=(--settings "/etc/node-red/config.js")

cd /opt || bashio::exit.nok "Could not change directory to Node-RED"

bashio::log.info "Starting Node-RED..."

exec npm start -- "${options[@]}"
