# node-red/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Runs the Nginx daemon
# ==============================================================================

# Wait for Node-RED to become available
bashio::net.wait_for 46836 localhost 300

bashio::log.info "Starting NGinx..."
exec nginx
